# blueproximity-plus/etc/xdg/autostart/blueproximity-plus-autostart.desktop
[Desktop Entry]
Exec=blueproximity-plus 2>> /tmp/bp_err.log
Name=Blueproximity-plus
Terminal=false
Type=Application
StartupNotify=false
